# deploy/alpha-factory.service
[Unit]
Description=Alpha Factory Async Engine (Production)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/huyhoang/alpha-factory-v2
User=huyhoang
Group=huyhoang

# Run directly with bash
ExecStart=/bin/bash -lc "cd /home/huyhoang/alpha-factory-v2 && ./scripts/run_vps.sh"

# Restart policy  
Restart=always
RestartSec=10

# Graceful shutdown
TimeoutStartSec=120
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
